# config/Debian/daemon/systemd/imscp_daemon.service
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
[Unit]
Description=i-MSCP communication daemon
After=local-fs.target network.target syslog.target

[Service]
Nice=10
Type=forking
PIDFile=${RUNSTATEDIR}/imscp_daemon.pid
ExecStart=${BIN_DIR}imscp_daemon -b ${SBIN_DIR}/imscp-rqst-mngr -p ${RUNSTATEDIR}/imscp/imscp_daemon.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
