# configs/Debian/systemd/bind9-resolvconf.service.d/imscp-fix-IP-1333.conf
# SYSTEMD.SERVICE(5) configuration file - auto-generated by i-MSCP
#     DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
# Drop-In file that provides a fix for #IP-1333
# See also: https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=744304
[Service]
RemainAfterExit=yes
